# deploy/snake-attack.service
[Unit]
Description=Ultimate Snake Attack — multiplayer game server
After=network.target

[Service]
Type=simple
User=snake
Group=snake
WorkingDirectory=/opt/snake-attack
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production

# Логи
StandardOutput=journal
StandardError=journal
SyslogIdentifier=snake-attack

[Install]
WantedBy=multi-user.target
